# gemini-api.service
[Unit]
Description=Gemini Web2API Proxy Server
After=network.target

[Service]
User={{USER}}
WorkingDirectory={{WORKDIR}}
ExecStart={{WORKDIR}}/.venv/bin/python gemini.py server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gemini-api

[Install]
WantedBy=multi-user.target
